# mcp-memory/mcp-memory.service
# L9 MCP Memory Server - Systemd Service File
# 
# Install:
#   sudo cp mcp-memory/mcp-memory.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable mcp-memory
#   sudo systemctl start mcp-memory
#
# Check status:
#   sudo systemctl status mcp-memory
#   sudo journalctl -u mcp-memory -f

[Unit]
Description=L9 MCP Memory Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/l9

# Environment (override in /etc/systemd/system/mcp-memory.service.d/override.conf)
EnvironmentFile=-/opt/l9/mcp-memory/.env

# Start command
ExecStart=/opt/l9/.venv/bin/uvicorn mcp-memory.main:app --host 127.0.0.1 --port 9001

# Restart policy
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/l9/logs
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-memory

[Install]
WantedBy=multi-user.target

